# roles/subspace_operator/templates/subspace_operator.service.j2

[Unit]
Description="subspace node"
After=network-online.target

[Service]
User={{ subspace_operator_user }}
Group={{ subspace_operator_group }}
WorkingDirectory={{ subspace_operator_root_path }}/
ExecStart={{ subspace_operator_binary_name }} \
    --chain {{ subspace_operator_chain_id }} \
    --name {{ ansible_hostname }} \
    --base-path {{ subspace_operator_node_path }}/ \
    --blocks-pruning 256 \
    --state-pruning archive \
    --prometheus-external \
    --dsn-listen-on /ip4/0.0.0.0/udp/{{ subspace_operator_dsn_port }}/quic-v1 \
    --dsn-listen-on /ip4/0.0.0.0/tcp/{{ subspace_operator_dsn_port }} \
    -- --domain-id 1 \
    --no-private-ipv4 \
    --operator-id 143 \
    --bootnodes /ip4/3.87.28.170/tcp/40333/p2p/12D3KooWGHtULvhdKMZtzigSK1438uWXPj9rBQHVzTaKMWv1WRXp \
    --keystore-path /keystore \
    --listen-addr /ip4/0.0.0.0/tcp/{{ subspace_operator_p2p_port }}


StandardOutput=append:{{ subspace_operator_log_folder }}/{{ subspace_operator_log_file }}
StandardError=append:{{ subspace_operator_log_folder }}/{{ subspace_operator_log_file }}


Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
